# config/spug.ini
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command = nginx -g "daemon off;"
autostart = true
priority = 10
stdout_logfile = /var/log/supervisor/nginx.log
stderr_logfile = /var/log/supervisor/nginx_error.log
redirect_stderr = false
autorestart = true
startretries = 3

[program:spug-api]
command = bash /data/spug/spug_api/tools/start-api.sh
directory = /data/spug/spug_api
autostart = true
startsecs = 5
priority = 20
stdout_logfile = /data/spug/spug_api/logs/api.log
stderr_logfile = /data/spug/spug_api/logs/api_error.log
redirect_stderr = false
autorestart = true
startretries = 5
stopwaitsecs = 10

[program:spug-ws]
command = bash /data/spug/spug_api/tools/start-ws.sh
directory = /data/spug/spug_api
autostart = true
startsecs = 5
priority = 30
stdout_logfile = /data/spug/spug_api/logs/ws.log
stderr_logfile = /data/spug/spug_api/logs/ws_error.log
redirect_stderr = false
autorestart = true
startretries = 5
stopwaitsecs = 10

[program:spug-worker]
command = bash /data/spug/spug_api/tools/start-worker.sh
directory = /data/spug/spug_api
autostart = true
startsecs = 5
priority = 40
stdout_logfile = /data/spug/spug_api/logs/worker.log
stderr_logfile = /data/spug/spug_api/logs/worker_error.log
redirect_stderr = false
autorestart = true
startretries = 5
stopwaitsecs = 10

[program:spug-monitor]
command = bash /data/spug/spug_api/tools/start-monitor.sh
directory = /data/spug/spug_api
autostart = true
startsecs = 5
priority = 50
stdout_logfile = /data/spug/spug_api/logs/monitor.log
stderr_logfile = /data/spug/spug_api/logs/monitor_error.log
redirect_stderr = false
autorestart = true
startretries = 3
stopwaitsecs = 10

[program:spug-scheduler]
command = bash /data/spug/spug_api/tools/start-scheduler.sh
directory = /data/spug/spug_api
autostart = true
startsecs = 5
priority = 60
stdout_logfile = /data/spug/spug_api/logs/scheduler.log
stderr_logfile = /data/spug/spug_api/logs/scheduler_error.log
redirect_stderr = false
autorestart = true
startretries = 3
stopwaitsecs = 10
